# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Exec=/home/taki/.config/conky/startconky.sh
Name=conky
